# Hugo/Linux/NVIDIA/TD-9-11-22/dump-database.service
[Unit]
Description=Execute dump-database.sh

[Service]
Type=oneshot
ExecStart=/home/hugo/dump-database.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target