# Hugo/Linux/RenduDu23/11/2022/dump-database.service
[Unit]
Description=Execute dump-database.sh

[Service]
Type=oneshot
ExecStart=/home/hugo/dump-database.sh
Remai,AfterExit=yes

[Install]
WantedBy=multi-user.target