# Eva/Linux/NVIDIA/TD-9-11-22/LabauvieRaffaelli/dump-database.service
[Unit]
Description=dump-database, save the data

[Service]
WorkingDirectory=/home/eva/sauvegarde
User=eva
ExecStart=/home/eva/dump-database.sh
Type=oneshot
