# Eva/Linux/NVIDIA/TD-9-11-22/LabauvieRaffaelli/dump-database.timer
[Unit]
Description= Run backup every 10 minutes

[Timer]
OnCalendar=*:0/10
AccuracySec=1s

[Install]
WantedBy=timers.targer
